# ops/native/ky-agent-executor-runtime-broker.socket
[Unit]
Description=AiCRM Agent Executor Runtime Broker Socket

[Socket]
ListenSequentialPacket=/run/aicrm-agent-runtime.sock
SocketUser=root
SocketGroup=ky-agent-executor
SocketMode=0660
RemoveOnStop=true

[Install]
WantedBy=sockets.target
