# ops/native/ky-agent-executor-runtime-broker.service
[Unit]
Description=AiCRM Agent Executor Root Runtime Broker
Requires=ky-agent-executor-runtime-broker.socket
After=ky-agent-executor-runtime-broker.socket

[Service]
Type=simple
User=root
Group=root
ExecStartPre=/usr/bin/install -d -o root -g ky-agent-executor -m 1770 /var/lib/aicrm-agent-executors
ExecStartPre=/usr/bin/install -d -o root -g root -m 0700 /var/lib/private/aicrm-codex-runtime
ExecStart=/data/kyai_crm/bin/ky-agent-executor-runtime-broker
Environment=KY_AGENT_EXECUTOR_CREDENTIAL_ROOT=/var/lib/aicrm-agent-executors
Environment=KY_AGENT_EXECUTOR_USER=ky-agent-executor
Environment=KY_CODEX_BINARY=/usr/bin/codex
Environment=KY_SYSTEMD_RUN_PATH=/usr/bin/systemd-run
Environment=KY_SYSTEMCTL_PATH=/usr/bin/systemctl
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
KillMode=control-group
UMask=0077
LimitCORE=0

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/aicrm-agent-executors /var/lib/private/aicrm-codex-runtime
InaccessiblePaths=/data/kyai_crm/config /data/kyai_crm/backups
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
# systemd 255's RestrictSUIDSGID seccomp profile rejects openat2; the broker
# requires openat2 for anchored no-symlink traversal and keeps NoNewPrivileges.
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_CHOWN CAP_FOWNER CAP_DAC_READ_SEARCH

[Install]
WantedBy=multi-user.target
