# ops/native/ky-agent-executor-runtime-broker.service
[Unit]
Description=AiCRM Agent Executor Root Runtime Broker
Requires=ky-agent-executor-runtime-broker.socket
After=ky-agent-executor-runtime-broker.socket

[Service]
Type=simple
User=root
Group=root
ExecStart=/data/kyai_crm/bin/ky-agent-executor-runtime-broker
Environment=KY_AGENT_EXECUTOR_CREDENTIAL_ROOT=/var/lib/aicrm-agent-executors
Environment=KY_AGENT_EXECUTOR_USER=ky-agent-executor
Environment=KY_CODEX_BINARY=/usr/bin/codex
Environment=KY_SYSTEMD_RUN_PATH=/usr/bin/systemd-run
Environment=KY_SYSTEMCTL_PATH=/usr/bin/systemctl
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
KillMode=control-group
UMask=0077
LimitCORE=0

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER
AmbientCapabilities=CAP_CHOWN CAP_FOWNER

[Install]
WantedBy=multi-user.target
