# ops/native/ky-agent-executor-service.service
[Unit]
Description=AiCRM Agent Executor Shadow Control Plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ky-agent-executor
Group=ky-agent-executor
WorkingDirectory=/data/kyai_crm
EnvironmentFile=/data/kyai_crm/config/ky-agent-executor-service.env
ExecStart=/data/kyai_crm/bin/ky-agent-executor-service
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
KillMode=control-group
UMask=0077
LimitCORE=0

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
